# packaging/scripts/prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
    systemctl stop effecty.service || true
    systemctl disable effecty.service || true
fi

#DEBHELPER#
